# files/etc/hotplug.d/iface/99-dslite-mtu
[ "$ACTION" = "ifup" ] || exit 0
[ "$INTERFACE" = "wan6_4" ] || exit 0
ip link set dev ds-wan6_4 mtu 1424
